# containers/auth/assets/s6-rc.d/pdns-run-prepare/up
if { mkdir -p /run/pdns }
if { chown nobody:nogroup /run/pdns }
chmod 02755 /run/pdns