# containers/auth/assets/s6-rc.d/pdns-log-prepare/up
if { mkdir -p /var/log/pdns }
if { chown nobody:nogroup /var/log/pdns }
chmod 02755 /var/log/pdns