# containers/auth/assets/s6-rc.d/pdns/run
#!/command/with-contenv /bin/sh
/bin/pdns-prepare-config
exec 2>&1
exec pdns_server