# containers/auth/assets/s6-rc.d/pdns-lib-prepare/up
if { mkdir -p /var/lib/pdns }
if { chown nobody:nogroup /var/lib/pdns }
chmod 02755 /var/lib/pdns