# containers/recursor/assets/s6-rc.d/pdns/run
#!/command/with-contenv /bin/sh
exec 2>&1
exec pdns_recursor